# cluelyguard.service
[Unit]
Description=CluelyGuard Anti-LLM Proctoring Daemon
Documentation=https://github.com/yourusername/cluelyguard
After=network.target pulseaudio.service
Wants=network.target

[Service]
Type=notify
User=cluelyguard
Group=cluelyguard
WorkingDirectory=/opt/cluelyguard
ExecStart=/usr/local/bin/cluelyguard-daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10



# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment variables
Environment=CLUELYGUARD_CONFIG=/etc/cluelyguard/local.yaml
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cluelyguard

[Install]
WantedBy=multi-user.target 